# images/apache/supervisord.conf
[supervisord]
nodaemon=true

[program:start]
command=/bin/bash -c "cd $PROJHOME && bin/start && cat *.cronjob | crontab"
autorestart=false
umask=0002

[program:cron]
command=/usr/sbin/cron -f
autorestart=true

[program:rsyslogd]
command=/usr/sbin/rsyslogd -d
autorestart=true

[program:apache2]
command=apache2-foreground
autorestart=true
umask=0002
