# images/sshgrid/supervisord.conf
[supervisord]
nodaemon=true

[program:rsyslogd]
command=/usr/sbin/rsyslogd -d
autorestart=true

[program:apache2]
command=/usr/sbin/sshd -d
autorestart=true
